# auto_updatedb_for_locate_command/updatedb.timer
[Unit]
Description=Run updatedb at boot

[Timer]
OnBootSec=5min
Unit=updatedb.service

[Install]
WantedBy=timers.target
